# crontab
0 19 * * 1-5 /usr/local/bin/python -u /app/src/main.py 2>&1 |tee -a /app/logs/stock.log
