# crontab
0 19 * * 1-5 python -u /app/src/main.py 2>&1 | tee /app/logs/stock.log
